# ipFireVariant/cron.daily
0 23 * * * /bin/sh -c 'd=$(date +\%Y-\%m-\%d); gzip -c /var/log/syslog-email-service.log > /var/log/syslog-email-service-$d.gz && : > /var/log/syslog-email-service.log'
0 23 * * * /bin/sh -c 'd=$(date +\%Y-\%m-\%d); gzip -c /var/log/syslog-email-service.log > /var/log/sysmail-$d.gz && : > /var/log/sysmail.log'
